# launchd/com.followupmonitor.ingestion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.followupmonitor.ingestion</string>

		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>/Users/sujalpatel/Documents/AI EMAIL SUMMARY /scripts/run-ingestion-once.sh</string>
		</array>

		<!-- Disabled by default: backend Nest cron already runs live sync every 2 minutes. -->
		<key>Disabled</key>
		<true/>

		<key>RunAtLoad</key>
		<false/>

    <!-- Every 2 minutes -->
    <key>StartInterval</key>
    <integer>120</integer>

    <key>StandardOutPath</key>
    <string>/Users/sujalpatel/Library/Logs/followupmonitor-ingestion.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sujalpatel/Library/Logs/followupmonitor-ingestion.err.log</string>
  </dict>
</plist>
