# launchd/com.followupmonitor.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.followupmonitor.backend</string>

    <key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>/Users/sujalpatel/Library/Application Support/followupmonitor/scripts/start-backend.sh</string>
		</array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>FOLLOWUPMONITOR_REPO</key>
      <string>/Users/sujalpatel/Documents/AI EMAIL SUMMARY </string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sujalpatel/Library/Logs/followupmonitor-backend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sujalpatel/Library/Logs/followupmonitor-backend.err.log</string>
  </dict>
</plist>
